# upstart-init-statsdaemon.conf
# Initscript for statsdaemon
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn

script
    /usr/bin/statsdaemon &>> /var/log/statsdaemon.log
end script
